# crontab
* * * * * /bin/bash $PROJECT_LOCATION/init.sh >> /var/log/ip-dns-updater.log
