# crontab
* * * * * cd [DIR]/ip-dns-updater && /bin/bash init.sh >> /var/log/ip-dns-updater.log
